# systemd/ParanoidBackupOneDrive.service
[Unit]
Description=test job service

[Service]
Type=oneshot
ExecStart=/home/fabian/Shared/WorkerTest