# systemd/ParanoidBackupOneDrive.timer
[Unit]
Description=ParanoidOneDriveBackup Timer

[Timer]
#OnUnitActiveSec=60s
OnCalendar=*-*-* 4:00:00

[Install]
WantedBy=timers.target